# ansible/roles/tailscale_ssh/templates/tailscale-autoauth.service.j2
# Tailscale automatic authentication service
# Ensures Tailscale is authenticated on boot for tagged devices
# Managed by Ansible - do not edit manually

[Unit]
Description=Tailscale Auto-Authentication
After=tailscaled.service
Requires=tailscaled.service
ConditionPathExists=/etc/tailscale/auth.conf

[Service]
Type=oneshot
EnvironmentFile=/etc/tailscale/auth.conf
ExecStartPre=/bin/sleep 5
ExecStart=/bin/bash -c '\
  if tailscale status 2>&1 | grep -q "Logged out"; then \
    echo "Tailscale logged out, re-authenticating..."; \
    tailscale up \
      --auth-key=${TAILSCALE_AUTH_KEY} \
      --ssh \
      --advertise-tags=${TAILSCALE_TAGS} \
      ${TAILSCALE_ACCEPT_ROUTES:+--accept-routes} \
      ${TAILSCALE_ADVERTISE_ROUTES:+--advertise-routes=${TAILSCALE_ADVERTISE_ROUTES}} \
      --accept-risk=lose-ssh \
      --reset; \
  else \
    echo "Tailscale already authenticated"; \
  fi'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
