# ansible/roles/tailscale_ssh/templates/tailscale-autoauth.service.j2
# Tailscale automatic authentication service
# Ensures Tailscale is authenticated on boot for tagged devices
# Managed by Ansible - do not edit manually

[Unit]
Description=Tailscale Auto-Authentication
After=tailscaled.service
Requires=tailscaled.service
ConditionPathExists=/etc/tailscale/auth.conf

[Service]
Type=oneshot
EnvironmentFile=/etc/tailscale/auth.conf
ExecStartPre=/bin/sleep 5
ExecStart=/bin/bash -c '\
  STATUS_JSON=$(tailscale status --json 2>/dev/null); \
  if echo "$STATUS_JSON" | grep -q "\"BackendState\":\"Running\""; then \
    echo "Tailscale already authenticated and running"; \
  else \
    echo "Tailscale needs authentication, re-authenticating..."; \
    tailscale up \
      --authkey=${TAILSCALE_AUTH_KEY} \
      --ssh \
      --advertise-tags=${TAILSCALE_TAGS} \
      $([ "$TAILSCALE_ACCEPT_ROUTES" = "1" ] && echo "--accept-routes" || echo "--accept-routes=false") \
      ${TAILSCALE_ADVERTISE_ROUTES:+--advertise-routes=${TAILSCALE_ADVERTISE_ROUTES}} \
      --accept-risk=lose-ssh \
      --reset; \
  fi'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
