# ansible/roles/cloudflared/templates/cloudflared.service.j2
[Unit]
Description=Cloudflare Tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel run --token {{ cloudflared_tunnel_token }}
Restart=always
RestartSec=5
TimeoutStartSec=0

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
